# qentra-infra-agent.service
[Unit]
Description=Qentra Infrastructure Agent
Documentation=https://github.com/Qentra-on-call/Qentral-Infrastructure-agent
After=network-online.target pveproxy.service
Wants=network-online.target
# Never permanently give up retrying. An earlier version of this unit had a
# MemoryMax tight enough to risk an OOM-kill loop; if that (or anything else)
# ever crash-loops the agent, systemd's default restart-rate limit would stop
# trying after 5 failures in 10s and the node would silently go dark. Disabling
# the limit means it always keeps retrying instead of failing permanently.
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=/etc/qentra-infra-agent/env
ExecStart=/usr/bin/node /opt/qentra-infra-agent/agent.js
Restart=always
RestartSec=5
# Needs root: pvesh talks to the local Proxmox API over a root-only unix
# socket, and zpool status requires root on most distros. Deliberately NOT
# sandboxed with ProtectSystem=strict — pvesh/pve-cluster's runtime sockets
# live under /run in locations that vary by Proxmox version, and getting that
# wrong would silently break collection with no clear error.
User=root
NoNewPrivileges=true
# Hard resource caps — this is a monitoring sidecar on a production hypervisor,
# it must never be able to compete with VM workloads for CPU/RAM even if it
# misbehaves (a runaway loop, a memory leak). The agent's actual footprint is a
# few `pvesh get` calls every 30s; these ceilings are generous headroom, not a
# tuned budget. MemoryMax was 128M — bumped to 256M after it looked like it may
# have been OOM-killing the process (Node's baseline RSS + `pvesh` child-process
# spawns can add up faster than expected under a tight cgroup cap).
# CPUQuota was 25% — confirmed live (a busy office-on-prem host with many
# VMs/pools) that a real collection cycle needs ~2-2.5min of actual CPU time
# to parse/process everything; at a 25% cap that FORCES ~8-10min of wall-clock
# per cycle, which collides almost exactly with the agent's own STALE_MIN
# (8min) stuck-ingest watchdog — the cap and the watchdog fight each other,
# guaranteeing a permanent restart loop that never ships data on any host
# whose real per-cycle CPU need exceeds roughly STALE_MIN x CPUQuota. Raised
# to 100%: the agent is still single-threaded (one JS thread + short-lived
# `pvesh`/`ipmitool` children), so this caps it at roughly ONE core's worth of
# usage regardless — negligible on any real multi-core hypervisor, and Nice=10
# + IOSchedulingClass=best-effort below are what actually protect VM workloads
# under contention, not this percentage.
CPUQuota=100%
MemoryMax=256M
TasksMax=128
# Below-normal scheduling priority so it never contends with VM/storage I/O.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
